# cron/.local/cron/cron.daily/back_up_WebDAV.bk
#!/bin/bash
# File              : back_up_WebDAV.sh
# Author            : zhangzhao <zhangzhao@ihep.ac.cn>
# Date              : 27.05.2020
# Last Modified Date: 09.02.2021
# Last Modified By  : zhangzhao <zhangzhao@ihep.ac.cn>
NAME=backup_$(date "+%Y%m%d").tar.gz
cd /home/zhangzhao
tar -czf "$NAME" --exclude=./dotfiles/.git dotfiles
mv $NAME IHEP/backup
