# etc/systemd/systemd/system/home-zhangzhao-Mount.mount
[Unit]
Description=Mount sdb

[Mount]
What=/dev/disk/by-uuid/ec992520-ff93-402b-850a-9d569ee348de
Where=/home/zhangzhao/Mount
Type=ext4
Options=defaults,rw
TimeoutSec=15

[Install]
WantedBy=local-fs.target
