# etc/systemd/systemd/system/mnt.mount
[Unit]
Description=Mount sdb

[Mount]
What=/dev/sdb
Where=/mnt
Options=defaults
Type=ext4
TimeoutSec=15

[Install]
WantedBy=local-fs.target
